# .config/autostart/imwheel.desktop
[Desktop Entry]
Name=imwheel
Type=Application
Exec=imwheel
Terminal=false
